# aws-deploy/grid_bot.service
[Unit]
Description=Grid Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Binance/grid_trading_bot
ExecStart=/home/ubuntu/Binance/venv/bin/python /home/ubuntu/Binance/grid_trading_bot/main.py
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# 内存和CPU限制
MemoryHigh=768M
MemoryMax=900M
CPUQuota=80%

# 环境变量
Environment="PATH=/home/ubuntu/Binance/venv/bin:/usr/bin:/usr/local/bin"

[Install]
WantedBy=multi-user.target